# package/HalInsight/files/insightd.init
#!/bin/sh /etc/rc.common

INSIGHTD=/usr/sbin/insightd
START=50
STOP=50

start()
{
	$INSIGHTD -D /dev/ttyS0
}

stop()
{
	$INSIGHTD -e
}


