# firmware_1/rootfs/etc/init.d/insightd
#!/bin/sh /etc/rc.common

INSIGHTD=/sbin/insightd
START=50
STOP=50

start()
{
	$INSIGHTD -D /dev/ttyS0
}

stop()
{
	$INSIGHTD -e
}


