# firmware_2/rootfs/usr/lib/opkg/info/WeMo.list
/sbin/wemoApp
/sbin/web/setup.xml
/sbin/web/eventservice.xml
/sbin/startWemo.sh
/etc/init.d/wemo
/etc/ver.txt
/etc/icon.jpg
